# crontab/crawl_all.sh
#!/usr/bin/env bash
cd /home/apps/qichacha_spider
sh crawl_bjda.sh &
sh crawl_jd.sh &
